# make/target/target-initscripts/files/initmodules_hl101
#!/bin/sh
#
# Load all necessary st modules
# 
# chkconfig: S 9 0
#
  
   echo "init frontpanel"
   modprobe proton
   modprobe simu_button
   mknod -m 644 /dev/vfd c 147 0
   #mknod -m 644 /dev/hdmi_cec c 149 0
   modprobe e2_proc

   echo "init stmfb"
   modprobe stmcore-display-stx7109c3
   modprobe stmfb display0=720x576-32@60:12m:pal:yuv:yuv

   echo "LOADING MODULES"   
   echo "LOAD" >/dev/vfd

   /bin/ustslave /dev/st231-1 /lib/firmware/audio.elf
   echo "load v firmware"
   /bin/ustslave /dev/st231-0 /lib/firmware/video.elf

   echo "init embx"
   modprobe embxshell
   modprobe embxmailbox mailbox0=0xb9211000:136:set2 mailbox1=0xb9212000:0 ;
   modprobe embxshm mailbox0=shm:0:7:0x60000000:0:256:256:0:1024:0xa4000000:0x10000000 ;
   modprobe mme_host transport0=shm

   echo "init AVS"
   modprobe avs type=stv6418
   player=`cat /proc/stb/player/version | awk '{ print substr( $0, length($0) - 2, length($0) ) }'`
   echo "init player $player"

   modprobe mmelog
   modprobe ksound
   modprobe p2div64
   modprobe stm_v4l2
   modprobe pseudocard
   #modprobe fe-core demod=stb090x tuner=stb6110x
   #modprobe fe-core demod=stv090x tuner=sharp7306
   modprobe fe-core demod=stb0899 tuner=stb6100
   modprobe pti waitMS=20 videoMem=4096
   modprobe stmdvb
   modprobe stm_monitor
   modprobe player2 discardlateframe=0
   modprobe sth264pp
   modprobe stmalloc
   modprobe platform
   modprobe silencegen
   modprobe smartcard debug=100
   modprobe bpamem
   
   # Load module for swap in ram
   mem_by_cpu=$(($(grep MemTotal /proc/meminfo | awk ' { print $2 } ')*256))
   modprobe zram num_devices=1
   echo $mem_by_cpu > /sys/block/zram0/disksize

   if [ ! `grep 'hdmicec.enabled=false' /etc/enigma2/settings` ]; then
    modprobe cec
   fi
   echo "1" > /proc/sys/kernel/printk #for more log change level or  remake line
   cat /etc/videomode > /proc/stb/video/videomode
   showiframe -p /boot/bootlogo.mvi &
   stfbcontrol a 0

	stfbcontrol a 255
	
   # Enable swap in ram
   mkswap /dev/zram0
   swapon /dev/zram0
   
	echo " e2 "